# deploy/systemd/passkey-client.service
# ============================================================================
# Systemd Service Unit for FIDO2 Passkey Client (QA Environment)
# ============================================================================
# Installation:
#   sudo cp passkey-client.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable passkey-client
#   sudo systemctl start passkey-client
# ============================================================================

[Unit]
Description=FIDO2 Passkey Client - QA Environment
Documentation=https://github.com/JHyunJung/Passkey_Client
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=passkey
Group=passkey
WorkingDirectory=/opt/passkey-client

# Environment variables
Environment="NODE_ENV=production"
Environment="PORT=8003"

# Start command - Run QA environment
ExecStart=/usr/bin/npm run dev:qa

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/passkey-client

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=passkey-client

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
